# systemd/chigemotsu_daily_reboot.service
[Unit]
Description=Daily System Reboot

[Service]
Type=oneshot
ExecStart=/sbin/reboot
