# systemd/chigemotsu_daily_reboot.timer
[Unit]
Description=Timer for Daily System Reboot

[Timer]
OnCalendar=*-*-* 23:59:00

[Install]
WantedBy=timers.target
