# systemd/libcamerify_motion.service
[Unit]
Description=motion camera service
After=network.target

[Service]
Type=simple
# Set the user to run the service. Replace <USER> with your actual username.
User=<USER>
# (Optional) Set the group to run the service. Replace <GROUP> with your actual group name, or remove this line if not needed.
Group=<GROUP>
WorkingDirectory=/tmp
ExecStart=libcamerify motion
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
